# qtrobot/face_feature_detection.service
[Unit]
Description=Face feature detector used in MigrAVE

[Service]
ExecStart=/bin/bash -c "source /home/qtrobot/robot/autostart/qt_robot.inc && wait_for_network && prepare_ros_environment && source /home/qtrobot/catkin_build_ws/devel/setup.bash && roslaunch roslaunch migrave_face_feature_detector face_feature_detector.launch"
KillSignal=SIGKILL
User=qtrobot
Restart=always

[Install]
WantedBy=multi-user.target
